# ops/systemd/warroom-collector.service
[Unit]
Description=War Room OS v3 prospective collector
After=network-online.target
Wants=network-online.target

[Service]
Type=oneshot
User=warroom
Group=warroom
WorkingDirectory=/opt/warroom
Environment=WARROOM_ROOT=/opt/warroom
Environment=PYTHONPATH=/opt/warroom/src
ExecStart=/opt/warroom/.venv/bin/python /opt/warroom/scripts/warroom.py collect-all
NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=strict
ProtectHome=true
ReadWritePaths=/opt/warroom/runtime
