# ops/systemd/warroom-collector.timer
[Unit]
Description=Run War Room OS collector every five minutes

[Timer]
OnCalendar=*:0/5
Persistent=true
RandomizedDelaySec=20
Unit=warroom-collector.service

[Install]
WantedBy=timers.target
