# ops/systemd/warroom-api.service
[Unit]
Description=War Room OS v3 research dashboard
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=warroom
Group=warroom
WorkingDirectory=/opt/warroom
Environment=WARROOM_ROOT=/opt/warroom
Environment=PYTHONPATH=/opt/warroom/src
ExecStart=/opt/warroom/.venv/bin/python /opt/warroom/scripts/warroom.py serve --host 127.0.0.1 --port 8080
Restart=on-failure
RestartSec=5
NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=strict
ProtectHome=true
ReadWritePaths=/opt/warroom/runtime

[Install]
WantedBy=multi-user.target
